# app-misc/n8n/files/n8n.service
[Unit]
Description=n8n workflow automation platform
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=n8n
Group=n8n
WorkingDirectory=/var/lib/n8n
Environment="N8N_LISTEN_ADDRESS=127.0.0.1"
Environment="N8N_PORT=5678"
Environment="N8N_USER_FOLDER=/var/lib/n8n"
ExecStart=/usr/bin/n8n start
Restart=on-failure
RestartSec=5s
RuntimeDirectory=n8n
RuntimeDirectoryMode=0750
StateDirectory=n8n
StateDirectoryMode=0750

[Install]
WantedBy=multi-user.target
